# Manjaro_kde/supervisor/supervisor.d/sfgonboardingjobs.conf
[program:sfgonboardingjobs]
process_name=%(program_name)s
command=php /srv/http/sfg-onboarding/artisan queue:work --queue=high,default
autostart=true
autorestart=true
user=kevin-manjaro

