# Manjaro_kde/supervisor/supervisor.d/sfgonboarding.conf
[program:websockets]
command=php /srv/http/sfg-onboarding/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=kevin-manjaro